# cpuinfo.service
[Unit]
Description=cpuinfo as json service

[Service]
Type=simple
ExecStart=/usr/local/bin/cpuinfo_svc.py
KillMode=process
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target